# supervisord.conf
[program:node]
command = node src/server.js
directory = /home/dotcloud/current
